# deploy/mypet.service
[Unit]
Description=MyPet Django Application
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/mypet
Environment="PATH=/var/www/mypet/.venv/bin"
EnvironmentFile=/var/www/mypet/.env
ExecStart=/var/www/mypet/.venv/bin/gunicorn \
    --config /var/www/mypet/deploy/gunicorn.conf.py \
    mypet_project.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=5

StandardOutput=append:/var/log/mypet/gunicorn-stdout.log
StandardError=append:/var/log/mypet/gunicorn-stderr.log

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/mypet/media /var/log/mypet /tmp/mypet-uploads

[Install]
WantedBy=multi-user.target
